# builtin/app/go/data/common/dev-dep/upstart.conf.tpl
description "{{ name }} - Generated by Otto"

respawn
respawn limit 15 5

start on runlevel [2345]
stop on runlevel [06]

# We don't want to restart too fast for deps
post-stop exec sleep 5

script
  /usr/local/bin/{{ name }} >>/var/log/{{ name }}.log 2>&1
end script
